# packager/ukai.conf
description     "ukai"

start on runlevel [2345]
stop on runlevel [!2345]

console log

exec `su ukai -c "/usr/bin/ukai server --config-file /etc/ukai/gohan.yaml"`

respawn
respawn limit 10 90